# contrib/supervisord.ini
[program:opentracker]
command=/usr/bin/opentracker
autostart=true
